# extra/svn/fetch_changesets.cron
#!/bin/sh
key=`cat /etc/redmine.key`
curl "http://localhost/sys/fetch_changesets?key=$key" 

